# etc/init/gogo-varnishstat.conf
# Upstart script for gogo-varnishstat

description "gogo-varnishstat"

respawn
respawn limit 15 5

start on runlevel [2345] and varnish
stop on runlevel [06]

exec /usr/local/bin/gogo-varnishstat